# src/airootfs/etc/profile.d/mark_executable.sh
#!/bin/bash

main() {
    chmod 777 ~/Desktop/Install.desktop
}

main "$0"