# src/airootfs/etc/profile.d/start-services.sh
#!/bin/bash

NetworkManager
pulseaudio -D
gsettings set org.gnome.desktop.interface color-scheme 'prefer-dark'