# src/airootfs/archOwO_install/usr/lib/systemd/system/initial_setups.service
[Unit]
Description=Initial setups

[Service]
ExecStart=/archOwO/initial-setups.sh

[Install]
WantedBy=multi-user.target