# src/airootfs/etc/systemd/system/multi-user.target.wants/pulseaudio.service
[Unit]
Description=Pulse Audio

[Service]
ExecStart=/usr/bin/pulseaudio -D

[Install]
WantedBy=multi-user.target